# runit/caddy/run
#!/bin/sh
exec caddy --conf /etc/caddy/Caddyfile --agree=true
